# etc/gfarm-s3-router.service.in
[Unit]
Description = gfarm s3 router (gunicorn app)
After = network.target

[Service]
#Type = notify
User = @with_wsgi_user@
#Group = @with_wsgi_group@
#RuntimeDirectory = gunicorn
WorkingDirectory = /
ExecStart = @with_gunicorn@/bin/gunicorn --error-logfile - --threads @router_nthreads@ --workers @router_nworkers@ -b @with_router_bind@ --chdir @with_router_homedir@/ --capture-output --timeout 180 --reload router:app

#ExecReload = /bin/kill -s HUP $MAINPID
#KillMode = mixed
#TimeoutStopSec = 5
PrivateTmp = true

[Install]
WantedBy = multi-user.target
