# etc/gfarm-s3-webui.service.in
[Unit]
Description = Gfarm-S3 WebUI (gunicorn app)
After = network.target

[Service]
Type = simple
User = @with_gfarm_s3_user@
Group = @with_gfarm_s3_group@
RuntimeDirectory = gfarm-s3-webui
WorkingDirectory = /
ExecStart = @prefix@/libexec/gfarm-s3-webui

#ExecReload = /bin/kill -s HUP $MAINPID
#KillMode = mixed
#TimeoutStopSec = 5
PrivateTmp = true

[Install]
WantedBy = multi-user.target
