# etc/gfarm-s3-router.service.in
[Unit]
Description = Gfarm S3 router (gunicorn app)
After = network.target

[Service]
Type = simple
User = @GFARM_S3_USERNAME@
Group = @GFARM_S3_GROUPNAME@
RuntimeDirectory = gfarm-s3-router
WorkingDirectory = /
ExecStart = @prefix@/bin/gfarm-s3-router

#ExecReload = /bin/kill -s HUP $MAINPID
#KillMode = mixed
#TimeoutStopSec = 5
PrivateTmp = false

[Install]
WantedBy = multi-user.target
