# etc/gfarm-s3-router.service.in
[Unit]
Description = gfarm s3 router (gunicorn app)
After = network.target

[Service]
#Type = notify
User = @with_gfarm_s3_user@
#Group = @with_gfarm_s3_group@
#RuntimeDirectory = gunicorn
WorkingDirectory = /
Environment = threads=@router_nthreads@ workers=@router_nworkers@ bind=@with_router_addr@ chdir=@with_gfarm_s3_homedir@/ timeout=180

#ExecStart = @with_gunicorn@/bin/gunicorn --error-logfile - --threads $threads --workers $workers --bind $bind --chdir $chdir --capture-output --timeout $timeout --reload router:app
ExecStart = @with_gunicorn@/bin/gunicorn --error-logfile - --threads $threads --workers $workers --bind $bind --chdir $chdir --capture-output --timeout $timeout router:app

#ExecReload = /bin/kill -s HUP $MAINPID
#KillMode = mixed
#TimeoutStopSec = 5
PrivateTmp = true

[Install]
WantedBy = multi-user.target
