# etc/gunicorn.service.in
[Unit]
Description = gunicorn daemon
#Requires = gunicorn.socket
After = network.target

[Service]
#Type = notify
User = @with_wsgi_user@
#Group = @with_wsgi_group@
#RuntimeDirectory = gunicorn
WorkingDirectory = /
ExecStart = @with_gunicorn@/bin/gunicorn --error-logfile /tmp/gunicorn.err --log-file /tmp/gunicorn.log -w 2 -b 0.0.0.0:@with_wsgi_port@ --chdir @WEBUI_BINDIR@ gfarms3.wsgi:application
#ExecReload = /bin/kill -s HUP $MAINPID
#KillMode = mixed
#TimeoutStopSec = 5
PrivateTmp = true

[Install]
WantedBy = multi-user.target
