# etc/gunicorn.service.in
[Unit]
Description = gunicorn daemon
#Requires = gunicorn.socket
After = network.target

[Service]
#Type = notify
User = @WSGI_USER@
#Group = user2
#RuntimeDirectory = gunicorn
WorkingDirectory = /
ExecStart = @GUNICORN@ --error-logfile /tmp/gunicorn.err --log-file /tmp/gunicorn.log -w 2 -b 0.0.0.0:@WSGI_PORT@ --chdir @WEBUI_BINDIR@ gfarms3.wsgi:application
#ExecReload = /bin/kill -s HUP $MAINPID
#KillMode = mixed
#TimeoutStopSec = 5
PrivateTmp = true

[Install]
WantedBy = multi-user.target
